# recipes-connectivity/prometheus/files/mongodb_exporter.service
[Unit]
Description=Mongodb Exporter

[Service]
ExecStart=/usr/bin/mongodb_exporter

[Install]
WantedBy=multi-user.target
